# .config/systemd/user/wacom.service
[Unit]
Description=Configure my Wacom tablet
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/nicozhao1908/.config/wacom-config.sh

[Install]
WantedBy=graphical-session.target
